# docker_images/edge/supervisord.conf
[supervisord]
nodaemon=true

[program:vnc]
command=/usr/bin/vncserver :0 -geometry 1280x720 -depth 24 -rfbport 5900 -passwd /root/.vnc/passwd -fg -AlwaysShared -localhost=0
autorestart=true